# app/root/etc/services.d/openssh/run
#!/usr/bin/with-contenv sh

exec gosu root /usr/sbin/sshd -D
