# app/root/etc/services.d/apache2/run
#!/usr/bin/with-contenv sh

exec gosu root /usr/sbin/apachectl -DFOREGROUND
